# mail/services/dovecot/finish
#!/bin/execlineb -P
rm /run/dovecot/master.pid
